# charms/autopkgtest-dispatcher-operator/app/units/autopkgtest.target
[Unit]
Description=autopkgtest workers

[Install]
WantedBy=multi-user.target
