# charms/autopkgtest-janitor-operator/app/units/autopkgtest-cleanup-lxd.service.j2
[Unit]
Description=Clean stale LXD instalces

[Service]
Type=exec
User={{ user }}
WorkingDirectory=~
TimeoutStartSec=20m
ExecStart=/usr/local/bin/cleanup-lxd
