# charms/autopkgtest-janitor-operator/app/units/autopkgtest-cleanup-swift.timer
[Unit]
Description=Timer to clean old swift objects

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
