# charms/autopkgtest-janitor-operator/app/units/autopkgtest-cleanup-lxd.timer
[Unit]
Description=Timer to cleanup stale LXD instances

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
