# charms/autopkgtest-janitor-operator/app/units/autopkgtest-build-image@.timer
[Unit]
Description=Timer to build autopkgtest testbed image (%i)

[Timer]
OnCalendar=daily
RandomizedDelaySec=24h

[Install]
WantedBy=timers.target
