# charms/autopkgtest-janitor-operator/app/units/autopkgtest-build-vm@.timer
[Unit]
Description=Timer to build autopkgtest testbed VM (%i)

[Timer]
OnCalendar=daily
RandomizedDelaySec=24h
