# charms/autopkgtest-website-operator/app/units/amqp-status-collector.timer
[Unit]
Description=Read running test status from RabbitMQ (running.json)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
