# charms/autopkgtest-website-operator/app/units/autopkgtest-db-publisher.timer
[Unit]
Description=publish autopkgtest.db (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
