# charms/autopkgtest-website-operator/app/units/autopkgtest-db-backup.timer
[Unit]
Description=backup autopkgtest.db (timer)

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
