# charms/autopkgtest-website-operator/app/units/autopkgtest-db-backup.service.j2
[Unit]
Description=backup autopkgtest.db

[Service]
User={{ user }}
EnvironmentFile=-/etc/environment.d/proxy.conf
Type=oneshot
ExecStart={{ webcontrol }}/autopkgtest-db-backup
TimeoutSec=600
