# charms/autopkgtest-website-operator/app/units/autopkgtest-db-backup.service.j2
[Unit]
Description=Backup sql database

[Service]
Type=oneshot
User={{ user }}
EnvironmentFile=-/etc/environment.d/proxy.conf
ExecStart={{ webcontrol }}/autopkgtest-db-backup
