# charms/autopkgtest-website-operator/app/units/autopkgtest-db-backup.timer
[Unit]
Description=Backup sql database

[Timer]
OnCalendar=*-*-* 00:05:00

[Install]
WantedBy=timers.target
