# charms/autopkgtest-website-operator/app/units/autopkgtest-update-github-jobs.timer
[Unit]
Description=Update GitHub job status (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
