# charms/autopkgtest-website-operator/app/units/autopkgtest-stats.service.j2
[Unit]
Description=Collect & submit package stats

[Service]
Type=oneshot
User={{ user }}
ExecStart={{ webcontrol }}/autopkgtest-stats
TimeoutStartSec=10m
