# charms/autopkgtest-website-operator/app/units/autopkgtest-stats.timer
[Unit]
Description=Collect & submit package stats

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
