# charms/autopkgtest-website-operator/app/units/autopkgtest-queue-collector.timer
[Unit]
Description=Cache queue entries from AMQP (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
