# charms/autopkgtest-website-operator/app/units/autopkgtest-web.target
[Unit]
Description=autopkgtest web

[Install]
WantedBy=multi-user.target
