# charms/autopkgtest-website-operator/app/units/autopkgtest-running-collector.service.j2
[Unit]
Description=Read running test status from RabbitMQ (running.json)

[Service]
User={{ user }}
Type=oneshot
RuntimeDirectory=autopkgtest-running-collector
RuntimeDirectoryPreserve=yes
ExecStart={{ webcontrol }}/autopkgtest-running-collector
