# charms/autopkgtest-website-operator/app/units/autopkgtest-db-publisher.service.j2
[Unit]
Description=publish autopkgtest.db

[Service]
User={{ user }}
EnvironmentFile=-/etc/environment.d/proxy.conf
Type=oneshot
ExecStart={{ webcontrol }}/autopkgtest-db-publisher
TimeoutSec=300
