# charms/autopkgtest-website-operator/app/units/cache-amqp.service.j2
[Unit]
Description=Cache queue entries from AMQP

[Service]
Type=oneshot
TimeoutStartSec=600
User=nobody
StateDirectory=cache-amqp
ExecStart={{ webcontrol }}/cache-amqp
