# charms/autopkgtest-website-operator/app/units/cache-ppa-containers.timer
[Unit]
Description=collect ppa containers (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
