# charms/autopkgtest-website-operator/app/units/cache-amqp.timer
[Unit]
Description=Cache queue entries from AMQP (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=autopkgtest-web.target
