# charms/autopkgtest-website-operator/app/units/db-backup.service.j2
[Unit]
Description=Backup sql database

[Service]
Type=oneshot
User={{ user }}
ExecStart={{ webcontrol }}/db-backup
