# charms/autopkgtest-website-operator/app/units/download-results.service.j2
[Unit]
Description=Download results

[Service]
User={{ user }}
ExecStart={{ webcontrol }}/download-results
Restart=on-failure
RestartSec=1m

[Install]
WantedBy=autopkgtest-web.target
