# charms/autopkgtest-website-operator/app/units/indexed-packages.timer
[Unit]
Description=Index test results (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
