# charms/autopkgtest-website-operator/app/units/indexed-packages.service.j2
[Unit]
Description=Index test results

[Service]
Type=oneshot
User={{ user }}
ExecStart={{ webcontrol }}/indexed-packages
