# charms/autopkgtest-website-operator/app/units/publish-db.timer
[Unit]
Description=publish autopkgtest.db (timer)

[Timer]
OnCalendar=minutely

[Install]
WantedBy=autopkgtest-web.target
