# charms/autopkgtest-website-operator/app/units/publish-db.service.j2
[Unit]
Description=publish autopkgtest.db

[Service]
User=ubuntu
Type=oneshot
ExecStart={{ webcontrol }}/publish-db
TimeoutSec=300
