# charms/autopkgtest-website-operator/app/units/publish-db.service.j2
[Unit]
Description=publish autopkgtest.db

[Service]
User={{ user }}
EnvironmentFile=-/etc/environment.d/proxy.conf
Type=oneshot
ExecStart={{ webcontrol }}/publish-db
TimeoutSec=300
