# charms/dispatcher-operator/units/cloud-worker-maintenance.timer
[Unit]
Description=Cloud worker maintenance (timer)

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1h

[Install]
WantedBy=autopkgtest.target
