# charms/dispatcher-operator/units/pull-package-configs.timer
[Unit]
Description=Pull package configs repo (timer)
After=autopkgtest.target

[Timer]
OnCalendar=minutely

[Install]
WantedBy=autopkgtest.target
