# charms/janitor-operator/units/build-adt-vm@.timer
[Unit]
Description=Timer to build autopkgtest testbed VM (%i)
After=autopkgtest.target

[Timer]
OnBootSec=24h
OnUnitInactiveSec=24h
RandomizedDelaySec=6h
